# legacy/deploy/clide-web.service
[Unit]
Description=Clide Web Server (FastAPI + tmux)
After=network.target

[Service]
Type=simple
User=jpmschweitzer
Group=jpmschweitzer
WorkingDirectory=/mnt/media/Projects

ExecStart=/mnt/media/Projects/clide/.venv/bin/clide-web

# Environment
Environment="CLIDE_WEB_PORT=8888"
Environment="CLIDE_WEB_PROJECTS_DIR=/mnt/media/Projects"
Environment="CLIDE_WEB_CLIDE_BIN=/mnt/media/Projects/clide/.venv/bin/clide"
Environment="TERM=xterm-256color"
Environment="COLORTERM=truecolor"
Environment="PATH=/home/jpmschweitzer/.nvm/versions/node/v24.11.0/bin:/usr/local/bin:/mnt/media/Projects/clide/.venv/bin:/usr/bin:/bin"

# Process management
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=10
SendSIGKILL=yes

Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
